# templates/clamav-unofficial-sigs.service.j2
# {{ ansible_managed }}
[Unit]
Description=Clamav Unofficial Sigs Update service
Documentation=man:clamav-unofficial-sigs(8)

[Service]
ExecStart=/usr/local/sbin/clamav-unofficial-sigs.sh
